# platforms/debian/packageinfo/apcupsd.conffiles
/etc/apcctrl/apcctrl.conf
/etc/init.d/apcupsd
/etc/logrotate.d/apcupsd
/etc/apcctrl/hosts.conf
/etc/apcctrl/multimon.conf
/etc/apcctrl/apccontrol
/etc/apcctrl/changeme
/etc/apcctrl/commfailure
/etc/apcctrl/commok
/etc/apcctrl/offbattery
/etc/apcctrl/onbattery
