# platforms/debian/packageinfo/postrm
#!/bin/sh -e

if [ "$1" = "purge" ] ; then
    update-rc.d apcupsd remove >/dev/null
    rm -f /etc/apcctrl/apcctrl.conf
    rm -f /etc/apcctrl/apcupsd.status
    rm -f /etc/apcctrl/apcupsd.events
    rm -f /etc/apcctrl
fi
