# platforms/debian/examples/cron.daily
#! /bin/sh

# apcctrl	Cron script to rotate apcctrl log files daily.
#
# Written by Kevin Foss <kfoss@mint.net>.
# 5/31/99

cd /var/log
LOG=apcctrl.log

if [ -f $LOG ]; then
   savelog -g adm -m 644 -u root -c 7 $LOG >/dev/null
fi

/etc/init.d/apcctrl restart &>/dev/null
